# supervisor/conf.d/minecraft.conf
[program:minecraft]
directory=/data/
command=java -Xmx400M -Xms400M -jar Tekkit.jar nogui

